# com.bookiebenchmark.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bookiebenchmark.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.12/bin/python3</string>
        <string>/Users/williamshiber/Desktop/githubpage/bookiebenchmark/main.py</string>
        <string>morning</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/williamshiber/Desktop/githubpage/bookiebenchmark</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/williamshiber/Desktop/githubpage/bookiebenchmark/logs/daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williamshiber/Desktop/githubpage/bookiebenchmark/logs/daily_error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
